# debian/postinst
#!/bin/sh

set -e

systemctl enable lider.service
systemctl start lider.service
